# backend/systemd-units/perona-schedule.timer
[Unit]
Description=Run perona operators on schedule

[Timer]
OnBootSec=1min
OnUnitActiveSec=60sec

[Install]
WantedBy=timers.target
